# python/lib/systemd/loudml.service
[Unit]
Description=Bonsai App

[Service]
Type=simple
ExecStart=/usr/bin/loudmld
Restart=always

[Install]
WantedBy=multi-user.target
